# package/yig.service
[Unit]
Description=yig
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/yig
ExecStop=/usr/bin/kill $MAINPID
Restart=always

[Install]
WantedBy=multi-user.target
